# deploy/ch.sailingnomads.syndicator-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Syndicator watch daemon - launchd agent for macOS.

  Normally the watcher runs ONLY on the Ubuntu server; use this agent only
  if the server is down for a longer period. The pipeline lock prevents
  double processing, but one watcher is the intended setup.

  Install:
    cp deploy/ch.sailingnomads.syndicator-watch.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/ch.sailingnomads.syndicator-watch.plist
  Remove:
    launchctl unload ~/Library/LaunchAgents/ch.sailingnomads.syndicator-watch.plist

  Fill in OPENAI_API_KEY below before loading.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ch.sailingnomads.syndicator-watch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>syndicator</string>
        <string>watch</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/benno/git/syndicator</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENAI_API_KEY</key>
        <string>FILL_ME_IN</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/syndicator-watch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/syndicator-watch.log</string>
</dict>
</plist>
